# deploy/systemd/climate-web.service
[Unit]
Description=Climate web frontend (Next.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=climate
Group=climate
WorkingDirectory=/opt/climate/source/web
EnvironmentFile=/etc/climate/web.env
ExecStart=/usr/bin/env npm run start -- --hostname 127.0.0.1 --port 3000
Restart=always
RestartSec=3
TimeoutStartSec=60
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/climate/source/web/.next

[Install]
WantedBy=multi-user.target
