# deploy/systemd/climate-backend.service
[Unit]
Description=Climate API backend (FastAPI/Uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=climate
Group=climate
WorkingDirectory=/opt/climate/source
EnvironmentFile=/etc/climate/backend.env
ExecStart=/opt/climate/venv/bin/uvicorn climate_api.main:app --host ${API_HOST} --port ${API_PORT} --workers 2
Restart=always
RestartSec=3
TimeoutStartSec=30
TimeoutStopSec=30
KillSignal=SIGINT
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/climate/source/data

[Install]
WantedBy=multi-user.target
